# config/brmgr@.service.in
[Unit]
Description=brmgr
After=network.target
Before=lxc.service lxc@.service

[Service]
Type=simple
ExecStartPre=@bindir@/brmgr-pre %i
ExecStart=@bindir@/brmgr %i --pid-file
ExecStopPost=@bindir@/brmgr-post %i

[Install]
WantedBy=lxc.service lxc@.service
