# config/brmgr.upstart
description "LXC Network"

start on starting lxc or starting lxc-instance
stop on stopped lxc

pre-start exec /usr/local/lib/lxc-net/bridge-up
exec /usr/local/lib/lxc-net/dnsmasq --pid-file
post-start exec touch /run/lxc/network_up
post-stop script
	/usr/local/lib/lxc-net/bridge-down
	exec rm -f /run/lxc/network_up
end script
