# config/brmgr.upstart.in
description "brmgr: manage bridge devices"

start on file FILE="/etc/brmgr/*.conf" EVENT=create
stop on runlevel [016]

instance ${MATCH}

pre-start script
	conf="${MATCH##/etc/brmgr/}"
	conf="${conf%%.conf}"
	exec @bindir@/brmgr-pre "$conf"
end script

script
	conf="${MATCH##/etc/brmgr/}"
	conf="${conf%%.conf}"
	exec @bindir@/brmgr "$conf" --pid-file
end script

post-stop script
	conf="${MATCH##/etc/brmgr/}"
	conf="${conf%%.conf}"
	exec @bindir@/brmgr-post "$conf"
end script
